# python/por/rpc/etc/service/tesseract/run
#!/bin/bash
set -ex

# STARTING RPC SERVER
exec python3 /app/tesseract/app.py
